# deploy/systemd/piper-http.service
[Unit]
Description=VisionClip Piper HTTP TTS Server
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/visionclip-piper-http
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
